# root/etc/services.d/postgres/run
#!/usr/bin/with-contenv sh
gosu postgres postgres
